# unwind.service
[Unit]
Description=Unwind Bedside Device
After=network.target graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=subu
WorkingDirectory=/home/subu/unwind
ExecStart=/home/subu/unwind/start.sh
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/subu/.Xauthority

[Install]
WantedBy=graphical-session.target
